# rpi-5/lightpuck.service
[Unit]
Description=Claude Lightpuck LED Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/bcavs/projects/claude-light-puck/rpi-5/main.py
WorkingDirectory=/home/bcavs/projects/claude-light-puck/rpi-5
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
